# deploy/weibo-ai-bridge.service
[Unit]
Description=Weibo AI Bridge Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=weibo
WorkingDirectory=/opt/weibo-ai-bridge
ExecStart=/opt/weibo-ai-bridge/weibo-ai-bridge
Restart=on-failure
RestartSec=5
Environment=PATH=/home/weibo/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=CONFIG_PATH=/etc/weibo-ai-bridge/config.toml
EnvironmentFile=-/etc/weibo-ai-bridge/.env

[Install]
WantedBy=multi-user.target
